# systemd/ize-ribbon-web.service
[Unit]
Description=Ize Ribbon web server
After=network.target

[Service]
Type=simple
User=ize
WorkingDirectory=/home/ize/ize-ribbon
Environment=PYTHONUNBUFFERED=1
RuntimeDirectory=ize-ribbon
RuntimeDirectoryMode=0775
ExecStart=/usr/bin/python3 /home/ize/ize-ribbon/webserver.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
