# systemd/ize-ribbon-network-fallback.service
[Unit]
Description=Ize Ribbon Wi-Fi fallback access point
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/ize/ize-ribbon
ExecStart=/usr/bin/python3 /home/ize/ize-ribbon/network_fallback.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
